# conf/systemd/tornado_icinga2_collector.service.d/neteye.conf
[Unit]
PartOf=neteye.target
Before=neteye.target